# crontab
TZ=Asia/Shanghai
*/30 * * * * /usr/bin/python /app/upload.py >/proc/1/fd/1 2>/proc/1/fd/2
