# crontab
TZ=Asia/Shanghai
*/5 * * * * /usr/bin/python /app/demo.py >/proc/1/fd/1 2>/proc/1/fd/2
